# debian/postinst
#!/bin/sh
set -e
systemctl daemon-reload
echo "proxkube installed. Enable the daemon with: systemctl enable --now proxkube-daemon"
