# debian/prerm
#!/bin/sh
set -e
systemctl stop proxkube-daemon 2>/dev/null || true
systemctl disable proxkube-daemon 2>/dev/null || true
